# board/raspberrypi/fs_edgemap_org/etc/systemd/system/macsec.service
[Unit]
Description=macsec networking
Before=nfs-client.target
After=multi-user.target

[Service]
Type=oneshot
WorkingDirectory=/opt/nk-macsec
ExecStart=/opt/nk-macsec/macsec.sh
RemainAfterExit=true
