# board/raspberrypi/fs_commnode/etc/udev/rules.d/10-usbmount.rules
KERNEL=="sd*[!0-9]|sr*", ENV{ID_SERIAL}!="?*", SUBSYSTEMS=="usb", RUN+="/opt/nfs/rpi-scripts/usb-insert.sh"
