# board/raspberrypi/fs_commnode/etc/udev/rules.d/20-usbremove.rules
SUBSYSTEM=="usb", ACTION=="remove", RUN+="/opt/nfs/rpi-scripts/usb-disconnect.sh $devpath"

